# .config/systemd/user/bluetooth-applet.service
[Unit]
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

[Service]
ExecStart=blueman-applet
Restart=on-failure
