# .config/systemd/user/copyq.service
[Unit]
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

[Service]
ExecStart=copyq --start-server
Restart=on-failure
